# deploy/raspi/argos.service.example
[Unit]
Description=ARGOS monitoring service
After=network.target

[Service]
Type=simple
# Replace 'pi' with the appropriate user on your system.
User=pi
WorkingDirectory=/opt/argos

# Use the virtual environment's python if installed, otherwise fall back
ExecStart=/usr/bin/env python3 -m argos_app --mode hardware --config /opt/argos/config/argos.yaml
Restart=on-failure

[Install]
WantedBy=multi-user.target